# conf/supervisord.conf
[group:stahlnow.com]
programs=web.stahlnow.com

[program:web.stahlnow.com]
directory=/opt/stahlnow
command=gunicorn
    project.wsgi:application
    --pythonpath /opt/stahlnow/env
    -w 4
    --bind 0.0.0.0:8008
user=stahl
autostart=true
autorestart=true
redirect_stderr=true
stopasgroup=true

